# rootfs/services/unbound/run
#!/bin/bash

logger -p info "s6-supervise : spawning unbound process"
exec unbound &> /dev/stdout
